# deploy/m7-window/systemd/canarysting-adapter.service
[Unit]
Description=CanarySting M7 Envoy ext_proc adapter (sockops cookie bridge + enforce + canary seeding)
Documentation=https://github.com/henleda/canarysting
After=network-online.target canarysting-staged-range.service
Wants=network-online.target

[Service]
# Runs as root: attaches the sockops cookie bridge + the enforce cgroup programs
# and seeds the negative-space canaries. Talks to the engine over gRPC.
EnvironmentFile=/etc/canarysting/m7.env
ExecStart=/opt/canarysting/bin/envoy-adapter \
  -listen 127.0.0.1:50051 \
  -engine 127.0.0.1:50052 \
  -scope ${SCOPE} \
  -inline
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
